# 1.7.3/provisions/assets/units/docker.service.d/50-before-deps.conf
[Unit]
Before=kubelet.service

[Service]
CPUAccounting=yes
MemoryAccounting=yes
BlockIOAccounting=yes
